# docker/supervisord.conf
; SPDX-License-Identifier: MPL-2.0
; Copyright 2026 Clinton Bunch
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:postgres]
command=/usr/local/bin/postgres -D /var/lib/thaum/postgresql/data
user=postgres
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:thaum]
command=/app/docker/run_thaum.sh
user=thaum
; Ensure credentials env from entrypoint (CREDENTIALS_DIRECTORY after /run/secrets staging) reaches gunicorn.
environment=HOME="/home/thaum",CREDENTIALS_DIRECTORY="%(ENV_CREDENTIALS_DIRECTORY)s",THAUM_CREDS_DIR="%(ENV_THAUM_CREDS_DIR)s"
directory=/app
autostart=true
; Restart only on unexpected exits; normal stop (exit 0) stays down.
autorestart=unexpected
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
startretries=15
startsecs=1
